# etc/cron.d/ipchangerd
#
# cron-job ipchangerd
#

MAILTO=root

*/5 * * * *     nobody if [ -x /usr/bin/ipchangerd ]; then /usr/bin/ipchangerd; fi
